# repos/community-staging-x86_64/fcron.service
[Unit]
Description=Cron Daemon

[Service]
Type=forking
PIDFile=/run/fcron.pid
ExecStart=/usr/sbin/fcron

[Install]
WantedBy=multi-user.target
