# repos/community-staging-x86_64/ptunnel.service
[Unit]
Description=Ping Tunnel
After=syslog.target

[Service]
Type=forking
PIDFile=/var/run/ptunnel.pid
ExecStart=/usr/bin/ptunnel \
 -setuid nobody -setgid nobody -setcon \
 -syslog \
 -daemon /var/run/ptunnel.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target
