# repos/community-staging-x86_64/vsftpd-ssl.service
[Unit]
Description=vsftpd daemon (legacy implicit SSL)

[Service]
ExecStart=/usr/bin/vsftpd -oimplicit_ssl=yes
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
