# repos/community-staging-x86_64/zcashd.service
[Unit]
Description=Zcash daemon
After=network.target

[Service]
ExecStart=/usr/bin/zcashd

[Install]
WantedBy=default.target
