# goaurrpc/repos/community-testing-x86_64/goaurrpc.service
[Unit]
Description=goaurrpc daemon
Requires=network-online.target
After=network-online.target

[Service]
User=goaurrpc
Group=goaurrpc
Restart=on-failure
ExecStart=/usr/bin/goaurrpc -c /etc/goaurrpc.conf -v
TimeoutStopSec=0
KillSignal=SIGTERM

# Hardening options
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
PrivateDevices=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
