# mailman/repos/extra-x86_64/mailman-checkdbs.timer
[Unit]
Description=Daily mailman checkdbs

[Timer]
OnCalendar=8:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
