# filebin/repos/community-any/filebin-user-cron.timer
[Unit]
Description=Run filebin user cron hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
