# cdemu-daemon/repos/community-testing-i686/cdemu-daemon.service
[Unit]
Description=CDEmu Daemon

[Service]
Type=dbus
BusName=net.sf.cdemu.CDEMUD_Daemon
ExecStart=/usr/lib/cdemu-daemon/cdemu-daemon-system.sh

[Install]
WantedBy=multi-user.target
