# repos/community-testing-x86_64/fossil.socket
[Unit]
Description=Fossil socket

[Socket]
# TCP port to listen on
ListenStream=8055
Accept=yes

[Install]
WantedBy=sockets.target
