# repos/community-testing-x86_64/iperf-udp.service
[Unit]
Description=Iperf UDP Server
After=network.target

[Service]
ExecStart=/usr/bin/iperf -u -s -V

[Install]
WantedBy=multi-user.target
