# repos/community-testing-x86_64/kresus.service
[Unit]
Description=Kresus service
After=network.target

[Service]
Type=simple
User=kresus
Group=kresus
PermissionsStartOnly=true
WorkingDirectory=~
StateDirectory=kresus
Environment=NODE_ENV=production
ExecStart=/usr/bin/kresus -c /etc/webapps/kresus/config.ini
Restart=always
CapabilityBoundingSet=
NoNewPrivileges=true
#SecureBits=noroot-locked
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
PrivateUsers=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=yes
ProtectKernelLogs=true
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
LockPersonality=true
#Not compatible with NodeJS
#MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
