# repos/community-testing-x86_64/lcdd.service
[Unit]
Description=LCD display driver
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/LCDd -c /etc/LCDd.conf
ControlGroup=cpu:/

[Install]
WantedBy=multi-user.target
