# repos/community-testing-x86_64/yate.service
[Unit]
Description=Next-generation telephony engine

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/yate
StandardOutput=syslog
Restart=always
RestartSec=60
RuntimeDirectory=yate

[Install]
WantedBy=multi-user.target
