# ucarp/trunk/ucarp.service
[Unit]
Description=UCARP daemon

[Service]
ExecStart=/usr/sbin/ucarp

[Install]
WantedBy=multi-user.target
