# keycloak/trunk/keycloak.service
[Unit]
Description=Keycloak server
After=network.target

[Service]
ExecStart=/opt/keycloak/bin/standalone.sh

[Install]
WantedBy=multi-user.target
