# sslh/trunk/next/sslh-select.service
[Unit]
Description=SSL/SSH multiplexer (select mode)
Conflicts=sslh-fork.service sslh.socket
After=network.target

[Service]
ExecStart=/usr/bin/sslh-select -F
PIDFile=/run/sslh.pid

[Install]
WantedBy=multi-user.target
