# hiawatha/repos/community-staging-x86_64/lefh-renew.service
[Unit]
Description=Renew Let's Encrypt certificates for Hiawatha
Wants=network-online.target
After=network-online.target

[Service]
Environment=HOME=%h
Type=oneshot
ExecStart=/usr/bin/lefh renew restart
