# hiawatha/repos/community-staging-x86_64/letsencrypt-renew.service
[Unit]
Description=Renew Let's Encrypt certificates
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/root/letsencrypt/letsencrypt renew restart
