# dehydrated/repos/community-testing-any/dehydrated.service
[Unit]
Description=Check and renew ACME TLS certificates

[Service]
Type=oneshot
ExecStart=/usr/bin/dehydrated -c
