# foldingathome/trunk/foldingathome.service
[Unit]
Description=Folding@home distributed computing client
After=network.target

[Service]
Type=simple
User=fah
DynamicUser=yes
StateDirectory=fah
WorkingDirectory=/var/lib/fah
ExecStart=/usr/bin/FAHClient --config /etc/foldingathome/config.xml

[Install]
WantedBy=multi-user.target

