# ejabberd/trunk/ejabberd.service
[Unit]
Description=A distributed, fault-tolerant Jabber/XMPP server
After=network.target

[Service]
Type=oneshot
User=ejabberd
Group=jabber
LimitNOFILE=16000
RestartSec=5
ExecStart=/usr/sbin/ejabberdctl \
		--config /etc/ejabberd/ejabberd.cfg \
		--ctl-config /etc/ejabberd/ejabberdctl.cfg \
		--logs "/var/log/ejabberd" \
		--spool "/var/lib/ejabberd" start
ExecStop=/usr/sbin/ejabberdctl stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
