# ucarp/repos/community-i686/ucarp.service
[Unit]
Description=UCARP daemon

[Service]
EnvironmentFile=/etc/conf.d/ucarp
ExecStart=/usr/sbin/ucarp $DAEMON_OPTS

[Install]
WantedBy=multi-user.target
