# aiccu/trunk/aiccu.service
[Unit]
Description=SixXS Automatic IPv6 Connectivity Configuration Utility
After=network.target

[Service]
Type=forking
PIDFile=/var/run/aiccu.pid
ExecStart=/usr/sbin/aiccu start
ExecStop=/usr/sbin/aiccu stop

[Install]
WantedBy=multi-user.target
