# emby-server/repos/community-i686/emby-server.service
[Unit]
Description=Emby brings together your videos, music, photos, and live television.

[Service]
EnvironmentFile=/etc/conf.d/emby-server
WorkingDirectory=/usr/lib/emby-server
ExecStart=/usr/bin/emby-server
User=emby
Group=emby
Restart=always

[Unit]
After=network.target

[Install]
WantedBy=multi-user.target
