# tenshi/repos/community-any/tenshi.service
[Unit]
Description=real-time log monitor

[Service]
ExecStart=/usr/sbin/tenshi -c /etc/tenshi/tenshi.conf -f

[Install]
WantedBy=multi-user.target
