# keepalived/trunk/keepalived.service
[Unit]
Description=Keepalive Daemon

[Service]
ExecStart=/usr/sbin/keepalived --dont-fork

[Install]
WantedBy=multi-user.target
