# miredo/repos/community-x86_64/miredo.service
[Unit]
Description=Teredo IPv6 Tunneling Daemon
After=network.target

[Service]
Type=forking
PIDFile=/run/miredo.pid
ExecStart=/usr/bin/miredo -c /etc/miredo/miredo.conf -p /run/miredo.pid

[Install]
WantedBy=multi-user.target
