# nodm/trunk/nodm.service
[Unit]
Description=nodm display manager

[Service]
EnvironmentFile=/etc/nodm.conf
ExecStart=/usr/bin/nodm

[Install]
Alias=display-manager.service
