# languagetool/trunk/languagetool.service
[Unit]
Description=Languagetool HTTP server

[Service]
Type=simple
PIDFile=/run/language-tool.pid
User=nobody
Restart=always
ExecStart=/usr/bin/languagetool --http --allow-origin "*"
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -QUIT $MAINPID

[Install]
WantedBy=multi-user.target
