# tinc/repos/community-x86_64/tincd@.service
[Unit]
Description=Tinc daemon for net %i
After=network.target

[Service]
ExecStart=/usr/bin/tincd -D --pidfile /run/tinc.%i.pid -n %I
ExecStop=/usr/bin/tincd -n %I -k
ExecReload=/usr/bin/tincd -n %I -k HUP
Restart=on-failure

[Install]
WantedBy=multi-user.target
