# mailman/repos/community-x86_64/mailman-nightlygzip.service
[Unit]
Description=Compression of mailman archive files
Documentation=file:///usr/lib/mailman/cron/crontab.in

[Service]
Type=oneshot
ExecStart=/usr/bin/python2 -S /usr/lib/mailman/cron/nightly_gzip
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
