# pdnsd/trunk/pdnsd.service
[Unit]
Description=proxy name server
After=network.target

[Service]
ExecStart=/usr/sbin/pdnsd

[Install]
WantedBy=multi-user.target
