# drbd/repos/community-x86_64/drbd.service
[Unit]
Description=Distributed Replicated Block Device
After=systemd-modules-load.service network.target

[Service]
Type=oneshot
RemainAfterExit=yes

# load config
ExecStart=/usr/sbin/drbdadm adjust all
# user interruptible version of wait-connect all
ExecStart=/usr/sbin/drbdadm wait-con-int
# become primary if configured
ExecStart=/usr/sbin/drbdadm sh-b-pri all
# disconnect and detach all resources
ExecStop=/usr/sbin/drbdadm down all

[Install]
WantedBy=multi-user.target
