# openethereum/repos/community-x86_64/parity-user.service
[Unit]
Description=Parity Ethereum Daemon
After=network.target

[Service]
ExecStart=/usr/bin/parity
Restart=on-failure

# Specifies which signal to use when killing a service. Defaults to SIGTERM.
# SIGHUP gives parity time to exit cleanly before SIGKILL (default 90s)
KillSignal=SIGHUP

[Install]
WantedBy=default.target
