# hyperion/repos/community-staging-x86_64/hyperiond.service
[Unit]
Description=Hyperiond Systemd service

[Service]
Type=simple
User=hyperion
Group=hyperion
ExecStart=/usr/bin/hyperiond /etc/hyperion/hyperion.config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
