# mailman3/repos/community-staging-any/mailman3-digests.timer
[Unit]
Description=Daily GNU Mailing List Digests
Documentation=https://mailman.readthedocs.io/en/latest/src/mailman/commands/docs/digests.html

[Timer]
AccuracySec=1us
OnCalendar=08:00:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
