# trunk/backuppc.service
[Unit]
Description=BackupPC

[Service]
User=backuppc
ExecStart=/usr/bin/BackupPC
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
