# trunk/backuppc.socket
[Unit]
Description=BackupPC Socket

[Socket]
ListenStream=/run/backuppc/backuppc.sock

[Install]
WantedBy=sockets.target
