# trunk/exim@.service
[Unit]
Description=Exim Mail Daemon per-connection server

[Service]
ExecStart=-/usr/sbin/exim -bs
StandardInput=socket
StandardError=syslog
