# trunk/iperf.service
[Unit]
Description=Iperf Server

[Service]
ExecStart=/usr/bin/iperf -s

[Install]
WantedBy=multi-user.target
