# trunk/mailman-senddigests.service
[Unit]
Description=Send mailman digests
Documentation=file:///usr/lib/mailman/cron/crontab.in

[Service]
Type=oneshot
ExecStart=/usr/bin/python2 -S /usr/lib/mailman/cron/senddigests
Nice=19
IOSchedulingClass=best-effort
IOSchedulingPriority=7
