# trunk/openfire.service
[Unit]
Description=Openfire XMPP Server
After=network.target

[Service]
Type=simple
User=openfire
Group=openfire
PrivateTmp=true
EnvironmentFile=/etc/conf.d/openfire
ExecStart=/usr/bin/java -server $OPENFIRE_OPTS -DopenfireHome=/var/lib/openfire -Dopenfire.lib.dir=/usr/lib/openfire -classpath /usr/lib/openfire/startup.jar -jar /usr/lib/openfire/startup.jar
SuccessExitStatus=143 SIGKILL

[Install]
WantedBy=multi-user.target
