# trunk/vboxservice.service
[Unit]
Description=VirtualBox Guest Services
ConditionVirtualization=oracle

[Service]
ExecStart=/usr/bin/VBoxService -f

[Install]
WantedBy=multi-user.target
