# logstash/repos/community-testing-x86_64/logstash@.service
.include /usr/lib/systemd/system/logstash.service

[Unit]
Description=logstash (%I)

[Service]
Environment=LOGSTASH_CONF=/etc/logstash/%I.conf
