# dbmail/trunk/dbmail-timsieved.service
[Unit]
Description=DBMail Sieve Server
After=syslog.target network.target mysqld.service postgresql.service

[Service]
Type=forking
ExecStart=/usr/sbin/dbmail-timsieved

[Install]
WantedBy=multi-user.target
