# dbmail/trunk/dbmail-timsieved.service
[Unit]
Description=DBMail Sieve Server
After=syslog.target network.target mysqld.service postgresql.service

[Service]
Type=forking
PIDFile=/var/run/dbmail/dbmail-timsieved.pid
ExecStart=/usr/sbin/dbmail-timsieved -p /var/run/dbmail/dbmail-timsieved.pid

[Install]
WantedBy=multi-user.target
