# xe-guest-utilities/trunk/xe-linux-distribution.service
[Unit]
Description=xen linux distribution service
Before=xend.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/bin/xe-linux-distribution /var/cache/xe-linux-distribution

[Install]
WantedBy=multi-user.target
