# nss-pam-ldapd/repos/community-x86_64/nslcd.service
[Unit]
Description=Naming services LDAP client daemon.
After=syslog.target network.target

[Service]
ExecStart=/usr/bin/nslcd --nofork

[Install]
WantedBy=multi-user.target
