# mosquitto/repos/community-staging-x86_64/mosquitto.service
[Unit]
Description=Mosquitto MQTT Broker daemon
Requires=network.target

[Service]
Type=notify
NotifyAccess=main
User=mosquitto
ExecStart=/usr/bin/mosquitto -c /etc/mosquitto/mosquitto.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

NoNewPrivileges=yes
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
