# conntrack-tools/repos/community-x86_64/conntrackd.service
[Unit]
Description=Conntrack Daemon

[Service]
ExecStart=/usr/sbin/conntrackd

[Install]
WantedBy=multi-user.target
