# bird/trunk/bird.service
[Unit]
Description=BIRD routing daemon

[Service]
ExecStart=/usr/sbin/bird -d
ExecReload=/usr/sbin/birdc configure
ExecStop=/usr/sbin/birdc down

[Install]
WantedBy=multi-user.target
