# iperf/trunk/iperf-udp.service
[Unit]
Description=Iperf UDP Server

[Service]
ExecStart=/usr/bin/iperf -u -s -V

[Install]
WantedBy=multi-user.target
