# iperf/trunk/iperf-tcp.service
[Unit]
Description=Iperf TCP Server

[Service]
ExecStart=/usr/bin/iperf -s -V

[Install]
WantedBy=multi-user.target
