# fssos-nsvs/repos/community-x86_64/nsvsd.service
[Unit]
Description=Name Service Via Sockets Daemon

[Service]
ExecStart=/usr/bin/nsvsd -f /etc/nsvsd.conf

[Install]
WantedBy=multi-user.target
