# dingo/repos/community-i686/dingo.service
[Unit]
Description=Dingo Secure DNS resolver
Wants=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/bin/dingo

[Install]
WantedBy=multi-user.target
