# quagga/repos/community-i686/babeld.service
[Unit]
Description=Babel routing daemon
BindTo=zebra.service
After=network.target zebra.service
ConditionPathExists=/etc/quagga/babeld.conf

[Service]
ExecStart=/usr/sbin/babeld
Restart=on-abort

[Install]
WantedBy=multi-user.target
