# quagga/repos/community-i686/zebra.service
[Unit]
Description=GNU Zebra routing manager
After=network.target
ConditionPathExists=/etc/quagga/zebra.conf

[Service]
ExecStartPre=/usr/sbin/ip route flush proto zebra
ExecStart=/usr/sbin/zebra
Restart=on-abort

[Install]
WantedBy=multi-user.target
