# backuppc/repos/community-testing-x86_64/backuppc.service
[Unit]
Description=BackupPC

[Service]
User=backuppc
PIDFile=/run/backuppc/BackupPC.pid
ExecStart=/usr/share/backuppc/bin/BackupPC
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
