# hiawatha/trunk/letsencrypt-renew.service
[Unit]
Description=Renew Let's Encrypt certificates

[Service]
Type=oneshot
ExecStart=/root/letsencrypt/letsencrypt renew restart
