# exabgp/trunk/exabgp.service
[Unit]
Description=ExaBGP
After=network.target

[Service]
ExecStart=/usr/bin/exabgp /etc/exabgp.conf

[Install]
WantedBy=multi-user.target
