# ipset/trunk/ipset.service
[Unit]
Description=IP Sets
Before=iptables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/ipset start
ExecStop=/usr/lib/systemd/scripts/ipset stop

[Install]
WantedBy=multi-user.target
