# osiris/repos/community-staging-x86_64/osirismd.service
[Unit]
Description=A file integrity management system

[Service]
User=osiris
Type=forking
ExecStart=/usr/bin/osirismd

[Install]
WantedBy=multi-user.target
