# lxc-docker/trunk/docker.service
[Unit]
Description=Docker Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/docker -d

[Install]
WantedBy=multi-user.target
