# incron/repos/community-i686/incrond.service
[Unit]
Description=Inotify System Scheduler

[Service]
Type=forking
PIDFile=/run/incrond.pid
ExecStart=/usr/sbin/incrond

[Install]
WantedBy=multi-user.target
