# repos/community-staging-any/dehydrated.timer
[Unit]
Description=Check and renew ACME TLS certificates

[Timer]
OnCalendar=daily UTC

[Install]
WantedBy=timers.target
