# repos/community-staging-any/home-assistant.service
[Unit]
Description=Home assistant
After=network.target

[Service]
Restart=on-failure
User=hass
Group=hass

ExecStart=/usr/bin/hass --config /var/lib/hass/ --log-rotate-days 1
RestartForceExitStatus=100

[Install]
WantedBy=multi-user.target
