# repos/community-i686/dbmail-lmtpd.service
[Unit]
Description=DBMail LMTP Server
After=syslog.target network.target mysqld.service postgresql.service

[Service]
Type=forking
PIDFile=/run/dbmail/dbmail-lmtpd.pid
ExecStart=/usr/bin/dbmail-lmtpd -p /run/dbmail/dbmail-lmtpd.pid

[Install]
WantedBy=multi-user.target
