# repos/community-i686/glusterd.service
[Unit]
Description=Distributed filesystem daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/glusterd

[Install]
WantedBy=multi-user.target
