# repos/community-i686/mailman-gatenews.timer
[Unit]
Description=Periodic mailman gate_news

[Timer]
OnCalendar=*/5:00:00
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
