# repos/community-i686/netperf.service
[Unit]
Description=Netperf Benchmark Server

[Service]
ExecStart=/usr/bin/netserver -D

[Install]
WantedBy=multi-user.target
