# repos/community-i686/vsftpd.service
[Unit]
Description=vsftpd daemon

[Service]
ExecStart=/usr/sbin/vsftpd
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
