# wesnoth/repos/community-testing-x86_64/wesnothd.service
[Unit]
Description=Wesnoth Server Daemon
After=network.target

[Service]
ExecStart=/usr/bin/wesnothd

[Install]
WantedBy=multi-user.target
