# repos/community-testing-i686/iperf-tcp.service
[Unit]
Description=Iperf TCP Server
After=network.target

[Service]
ExecStart=/usr/bin/iperf -s -V

[Install]
WantedBy=multi-user.target
