# repos/community-testing-i686/nrpe.socket
[Unit]
Description=NRPE Socket for Per-Connection Servers

[Socket]
ListenStream=5666
Accept=yes

[Install]
WantedBy=sockets.target
