# nullmailer/repos/community-staging-x86_64/nullmailer.service
[Unit]
Description=Nullmailer relay-only MTA
After=local-fs.target
ConditionPathExists=/var/nullmailer/queue

[Install]
WantedBy=multi-user.target

[Service]
WorkingDirectory=/var/nullmailer
ExecStart=/usr/bin/nullmailer-send
User=nullmail
Restart=always
