# ipsec-tools/repos/community-x86_64/ipsec.service
[Unit]
Description=Load IPSec Security Policy Database
After=syslog.target network.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/sbin/setkey -f /etc/ipsec.conf
ExecStop=/usr/sbin/setkey -F -P ; /usr/sbin/setkey -F

[Install]
WantedBy=multi-user.target
