# repos/community-any/dehydrated.timer
[Unit]
Description=Check and renew ACME TLS certificates

[Timer]
OnCalendar=daily UTC
RandomizedDelaySec=8m

[Install]
WantedBy=timers.target
