# repos/community-any/ferm.service
[Unit]
Description=for Easy Rule Making
After=ipset.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/sbin/ferm /etc/ferm.conf
ExecReload=/usr/sbin/ferm /etc/ferm.conf
ExecStop=/usr/sbin/ferm -F /etc/ferm.conf

[Install]
WantedBy=multi-user.target
