# repos/community-any/jmc.service
[Unit]
Description=Jabber mail component

[Service]
Type=simple
PIDFile=/var/run/jmc.pid
ExecStart=python2 /usr/lib/jmc/jmc.py -c /etc/ejabberd/jmc.xml

[Install]
WantedBy=bluetooth.target
