# repos/community-any/opensearch-keystore@.service
[Unit]
Description=OpenSearch Keystore Generation
ConditionPathExists=|!/etc/opensearch/%I/opensearch.keystore

[Service]
Type=oneshot
Group=opensearch
UMask=0007
Environment=ES_PATH_CONF=/etc/opensearch/%I
ExecStart=/usr/share/opensearch/bin/opensearch-keystore create
RemainAfterExit=yes
