# sage-notebook/repos/community-any/sage.service
[Unit]
Description=A free open-source mathematics software system

[Service]
ExecStart=/usr/bin/sage -n

[Install]
WantedBy=default.target
