# drone-runner-exec/repos/community-x86_64/systemd.service
[Unit]
Description=Drone Runner (Exec)
Requires=network-online.target
After=network-online.target

[Service]
User=drone-runner-exec
Group=drone-runner-exec
WorkingDirectory=/var/lib/drone-runner-exec
ReadWritePaths=/var/lib/drone-runner-exec
ExecStart=/usr/bin/drone-runner-exec

[Install]
WantedBy=multi-user.target
