# darkhttpd/repos/community-i686/darkhttpd.socket
[Unit]
Conflicts=darkhttpd.service

[Socket]
ListenStream=80
Accept=no

[Install]
WantedBy=sockets.target
