# mfs/trunk/mfschunkserver.service
[Unit]
Description=Distributed file system (chunk server)
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/mfschunkserver start
ExecStop=/usr/bin/mfschunkserver stop
ExecReload=/usr/bin/mfschunkserver reload
PIDFile=/run/mfs/mfschunkserver.pid
Restart=always

[Install]
WantedBy=multi-user.target
