# parity-ethereum/repos/community-x86_64/parity.service
[Unit]
Description=Parity Daemon
After=network.target

[Service]
ExecStart=/usr/bin/parity

[Install]
WantedBy=default.target
