# quagga/repos/community-staging-x86_64/ospfd.service
[Unit]
Description=OSPF routing daemon
BindTo=zebra.service
After=network.target zebra.service
ConditionPathExists=/etc/quagga/ospfd.conf

[Service]
ExecStart=/usr/sbin/ospfd
Restart=on-abort

[Install]
WantedBy=multi-user.target
