# quagga/repos/community-staging-x86_64/ripngd.service
[Unit]
Description=RIP routing daemon for IPv6
BindTo=zebra.service
After=network.target zebra.service
ConditionPathExists=/etc/quagga/ripngd.conf

[Service]
ExecStart=/usr/sbin/ripngd
Restart=on-abort

[Install]
WantedBy=multi-user.target
