# dictd/trunk/dictd.service
[Unit]
Description=Dictd Dictionary Server Daemon
After=network.target

[Service]
Type=forking
EnvironmentFile=-/etc/default/dictd
ExecStart=/usr/bin/dictd $DICTD_ARGS -- $DICTD_EARGS

[Install]
WantedBy=multi-user.target
