# mailman/trunk/mailman-cullbadshunt.timer
[Unit]
Description=Daily mailman cull_bad_shunt

[Timer]
OnCalendar=04:00:00
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
