# mfs/repos/community-i686/mfscgiserv.service
[Unit]
Description=Distributed file system (web interface)
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/mfscgiserv
ExecStop=/bin/kill -15 $MAINPID
PIDFile=/run/mfs/mfscgiserv.pid
Restart=always

[Install]
WantedBy=multi-user.target
