# ucarp/repos/community-staging-x86_64/ucarp.service
[Unit]
Description=UCARP daemon

[Service]
ExecStart=/usr/bin/ucarp

[Install]
WantedBy=multi-user.target
