# hhvm/repos/community-testing-x86_64/hhvm@.service
[Unit]
Description=HHVM FastCGI application server
After=syslog.target network.target

[Service]
PIDFile=/run/hhvm/%I.pid
PrivateTmp=true
User=http
Group=http
ExecStart=/usr/bin/hhvm --config=/etc/hhvm/php.ini --config=/etc/hhvm/%I.ini --mode=daemon -vServer.Type=fastcgi -vPidFile=/run/hhvm/%I.pid
Restart=on-failure
RuntimeDirectory=hhvm

[Install]
WantedBy=multi-user.target
