# keepalived/repos/community-i686/keepalived.service
[Unit]
Description=Keepalive Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/keepalived --dont-fork

[Install]
WantedBy=multi-user.target
