# ulogd/repos/community-x86_64/ulogd.service
[Unit]
Description=Netfilter Ulogd daemon
After=network.target

[Service]
ExecStart=/usr/sbin/ulogd

[Install]
WantedBy=multi-user.target
