# kubernetes/repos/community-testing-x86_64/kube-proxy.service
[Unit]
Description=Kubernetes Kube-Proxy Server
Documentation=man:kube-proxy(1)
After=network.target

[Service]
EnvironmentFile=-/etc/kubernetes/kube-proxy.env
ExecStart=/usr/bin/kube-proxy $KUBE_PROXY_ARGS
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
