# goagent/repos/community-any/goagent.service
[Unit]
Description=GoAgent Service
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=/usr/share/goagent/local/proxy.py

[Install]
WantedBy=multi-user.target

