# rabbitmq/repos/community-any/rabbitmq.service
[Unit]
Description=RabbitMQ Messaging Server
After=network.target

[Service]
Type=simple
User=rabbitmq
SyslogIdentifier=rabbitmq
EnvironmentFile=/etc/rabbitmq/rabbitmq-env.conf

ExecStart=/usr/bin/rabbitmq-server
ExecStop=/usr/bin/rabbitmqctl stop

[Install]
WantedBy=multi-user.target
