# repos/community-testing-any/glances.service
[Unit]
Description=Glances Server

[Service]
ExecStart=/usr/bin/glances -s

[Install]
WantedBy=multi-user.target
