# repos/community-testing-any/mailman3.service
[Unit]
Description=GNU Mailing List Manager
Documentation=https://mailman.readthedocs.io/en/latest
After=network.target

[Service]
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/mailman start -f
ExecReload=/usr/bin/mailman restart
ExecStop=/usr/bin/mailman stop
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PIDFile=/run/mailman/master.pid
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
RemoveIPC=true
RestrictAddressFamilies=~AF_PACKET
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadOnlyPaths=/etc/mailman.cfg -/etc/mailman.d
ReadWritePaths=/var/lock/mailman /var/spool/mailman
RuntimeDirectory=mailman
StateDirectory=mailman
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@resources
LogsDirectory=mailman
Type=forking
User=mailman

[Install]
WantedBy=multi-user.target
