# repos/community-testing-any/mailman3-gatenews.timer
[Unit]
Description=Hourly GNU Mailing List Digests
Documentation=https://mailman.readthedocs.io/en/latest/src/mailman/commands/docs/gatenews.html

[Timer]
AccuracySec=1us
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
