# repos/community-testing-any/nextcloud-cron.service
[Unit]
Description=Nextcloud; run cron.php

[Service]
KillMode=mixed
User=nextcloud
ExecStart=/usr/bin/PHP_INTERPRETER -f /usr/share/webapps/nextcloud/cron.php
