# repos/community-testing-any/supervisord.service
[Unit]
Description=Process Monitoring and Control Daemon
Documentation=http://supervisord.org
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/supervisord -c /etc/supervisord.conf
ExecReload=/usr/bin/supervisorctl reload
ExecStop=/usr/bin/supervisorctl shutdown

[Install]
WantedBy=multi-user.target
