# audit/repos/community-i686/auditd.service
[Unit]
Description=Audit Daemon

[Service]
ExecStartPre=/sbin/auditctl -R /etc/audit/audit.rules
ExecStart=/sbin/auditd -n
ExecStop=/bin/kill -15 $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
