# quagga/repos/community-testing-x86_64/isisd.service
[Unit]
Description=IS-IS routing daemon
BindTo=zebra.service
After=network.target zebra.service
ConditionPathExists=/etc/quagga/isisd.conf

[Service]
ExecStart=/usr/sbin/isisd
Restart=on-abort

[Install]
WantedBy=multi-user.target
