# repos/community-x86_64/dopewars.socket
[Unit]
Description=Dopewars Server
PartOf=dopewars.service

[Socket]
ListenStream=7902

[Install]
WantedBy=sockets.target
