# repos/community-x86_64/filebeat.service
[Unit]
Description=Log file forwarder to remote sources
Documentation=https://www.elastic.co/guide/en/beats/filebeat/current/index.html
After=network.target

[Service]
User=filebeat
ExecStart=/usr/bin/filebeat -c /etc/filebeat/filebeat.yml -path.home /usr/share/filebeat -path.config /etc/filebeat -path.data /var/lib/filebeat -path.logs /var/log/filebeat
Restart=on-failure

[Install]
WantedBy=multi-user.target
