# repos/community-x86_64/geth-user.service
[Unit]
Description=Geth
After=network.target

[Service]
ExecStart=/usr/bin/geth
KillSignal=SIGINT
Restart=on-failure

[Install]
WantedBy=default.target
