# repos/community-x86_64/minidlna.service
[Unit]
Description=minidlna server
After=network.target

[Service]
Type=simple
User=nobody
Group=nobody
ExecStart=/usr/bin/minidlnad -S
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target
