# repos/community-x86_64/motion.service
[Unit]
Description=Motion daemon
After=local-fs.target

[Service]
ExecStart=/usr/bin/motion
Type=forking
#StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
