# repos/community-x86_64/webhook.service
[Unit]
Description=Webhook Service
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/webhook -hooks /etc/webhook/hooks.json -verbose
Restart=on-failure

[Install]
WantedBy=multi-user.target
