# pkgtools/trunk/pkgfile.cron
#!/bin/bash

/usr/bin/pkgfile --update >/dev/null
