# dkms/repos/community-testing-any/dkms.service
[Unit]
Description=Dynamic Kernel Modules System

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/lib/systemd/scripts/dkms start
ExecStop=/usr/lib/systemd/scripts/dkms stop

[Install]
WantedBy=multi-user.target
