# polipo/repos/community-staging-x86_64/polipo.cron
#!/bin/bash
/etc/rc.d/polipo purge >/dev/null 2>&1
