# solr/repos/community-testing-any/solr.service
[Unit]
Description=Solr full text search engine
After=network.target

[Service]
Type=simple
Environment=SOLR_PID_DIR=/run/solr/
Environment=SOLR_LOGS_DIR=/var/log/solr/
User=solr
WorkingDirectory=/var/lib/solr
ExecStart=/usr/bin/solr start -f -d /usr/share/solr/server -s /var/lib/solr
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
NoNewPrivileges=true
LockPersonality=true
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=~CAP_NET_BROADCAST CAP_NET_RAW
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
RemoveIPC=true
RestrictAddressFamilies=~AF_PACKET AF_NETLINK AF_UNIX
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
RuntimeDirectory=solr
StateDirectory=solr
LogsDirectory=solr

[Install]
WantedBy=multi-user.target
