# slmodem-utils/trunk/slmodem.service
[Unit]
Description=SLModem daemon

[Service]
Type=simple
ExecStart=/usr/bin/slmodemd --country=USA --alsa

[Install]
WantedBy=multi-user.target
