# fcron/repos/community-i686/fcron.service
[Unit]
Description=Cron Daemon

[Service]
ExecStart=/usr/sbin/fcron --foreground --nosyslog
KillMode=process

[Install]
WantedBy=multi-user.target
